# py/crontab
* * * * * python3 ~/src/solarbattman/py/solarbattman.py
